# scripts/crontab/monthly.fish
#!/usr/bin/env fish

clorange reboots reset
clorange updates reset