# scripts/crontab/saturday.fish
#!/usr/bin/env fish

kitty --detach -T task holup new music