# waybar/waybar-config.path
[Path]
PathModified=%h/fes/dot/waybar/waybar.jsonc

[Install]
WantedBy=paths.target
