# waybar/waybar-config.path
[Path]
PathChanged=%h/.config/waybar/waybar.jsonc

[Install]
WantedBy=paths.target
