# scripts/systemd/services/18.service
[Service]
ExecStart=/home/axlefublr/prog/dotfiles/scripts/systemd/executables/18.fish
Environment=DISPLAY=:0