# scripts/systemd/services/27.service
[Service]
ExecStart=/home/axlefublr/prog/dotfiles/scripts/systemd/executables/27.fish
Environment=DISPLAY=:0