# scripts/systemd/services/minute.service
[Service]
ExecStartPre=/home/axlefublr/r/dot/scripts/processwait.fish
ExecStart=/home/axlefublr/r/dot/scripts/systemd/executables/minute.fish