# scripts/systemd/services/monday.service
[Service]
ExecStart=/home/axlefublr/prog/dotfiles/scripts/systemd/executables/monday.fish