# scripts/systemd/services/monthly.service
[Service]
ExecStart=/home/axlefublr/prog/dotfiles/scripts/systemd/executables/monthly.fish
Environment=DISPLAY=:0