# scripts/systemd/services/saturday.service
[Service]
ExecStart=/home/axlefublr/prog/dotfiles/scripts/systemd/executables/saturday.fish
Environment=DISPLAY=:0