# scripts/systemd/services/sunday.service
[Service]
ExecStart=/home/axlefublr/prog/dotfiles/scripts/systemd/executables/sunday.fish
Environment=DISPLAY=:0