# scripts/systemd/services/thursday.service
[Service]
ExecStart=/home/axlefublr/prog/dotfiles/scripts/systemd/executables/thursday.fish
Environment=DISPLAY=:0