# systemd/axleizer.service
[Unit]
After=graphical-session.target
PartOf=graphical-session.target
Requisite=graphical-session.target

[Service]
Type=exec
ExecStart=:env dash -c 'TELOXIDE_TOKEN="$(cat ~/fes/uviw/lvia/axleizer)" axleizer'

[Install]
UpheldBy=graphical-session.target
