# systemd/daily.service
[Unit]
Wants=graphical-session.target
After=graphical-session.target

[Service]
Type=oneshot
ExecStart=/home/axlefublr/fes/dot/systemd/daily.fish
