# systemd/daily.service
[Service]
ExecStartPre=/home/axlefublr/r/dot/scripts/processwait.fish
ExecStart=/home/axlefublr/r/dot/systemd/daily.fish
