# systemd/download-drop.path
[Path]
PathChanged=%h/wlx

[Install]
WantedBy=paths.target
