# systemd/download-drop.path
[Path]
PathChanged=%h/wlx
PathChanged=%h/wlx/tabs

[Install]
WantedBy=paths.target
