# systemd/flipboard.service
[Service]
Type=exec
ExecStart=%h/r/dot/scripts/flipboard.fish

[Install]
WantedBy=graphical-session.target
