# systemd/flipboard.path
[Path]
PathChanged=%h/.local/share/flipboard

[Install]
WantedBy=paths.target
