# systemd/flipboard.path
[Path]
PathExists=/home/axlefublr/.local/share/flipboard

[Install]
WantedBy=paths.target
