# systemd/flipboard.service
[Service]
Type=oneshot
ExecStart=env fish -c 'cat ~/.local/share/flipboard | copy'
ExecStart=env truncate -s 0 %h/.local/share/flipboard
