# systemd/flipboard.service
[Service]
Type=oneshot
TimeoutStartSec=1
ExecStart=%h/fes/dot/eli/flipboard.fish
