# systemd/flipboard.service
[Service]
Type=oneshot
TimeoutStartSec=2
ExecStart=%h/r/dot/scripts/flipboard.fish
