# systemd/frizz.service
[Unit]
StartLimitBurst=0
StartLimitIntervalSec=0

[Service]
Type=oneshot
ExecStart=env fish -c frizz
