# systemd/frizz.path
[Path]
PathModified=%h/.local/share/magazine/C

[Install]
WantedBy=paths.target
