# systemd/mandb.service
[Service]
ExecStart=env fish -c "systemd_mandb"
Restart=always
Type=exec

[Install]
WantedBy=graphical-session.target
