# systemd/minute.service
[Unit]
Wants=graphical-session.target
After=graphical-session.target

[Service]
Type=oneshot
ExecStart=fish -c systemd_minute
