# systemd/ten-minutes.service
[Service]
ExecStartPre=/home/axlefublr/r/dot/scripts/processwait.fish
ExecStart=fish -c systemd_ten_minutes
