# systemd/wake.path
[Path]
PathExists=/tmp/cami-just-suspended

[Install]
WantedBy=paths.target
