# systemd/wake.path
[Path]
PathExists=/tmp/mine/just-suspended

[Install]
WantedBy=paths.target
