# kanata/kanata.service
[Unit]
Requisite=graphical-session.target

[Service]
